# scripts/bazzite-reboot-notify.service
[Unit]
Description=Send alert to alerts.itinyk.app on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/var/home/Og/scripts/bazzite-reboot-notify.sh

[Install]
WantedBy=default.target
