# scripts/steamdeck-auto-update.timer
[Unit]
Description=Run Steam Deck auto-update daily at 1am (catches up after sleep/power off)

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
