# configs/bazzite/proton-ge-update.timer
[Unit]
Description=Check for Proton GE updates weekly

[Timer]
OnCalendar=Mon *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
