# configs/bazzite/proton-ge-update.service
[Unit]
Description=Proton GE Auto-Updater

[Service]
Type=oneshot
ExecStart=/var/home/Og/scripts/proton-ge-update.sh
